# bin/crontab
# NOTE: master is ~steve/bin/crontab
8 3 * * * $HOME/bin/daily
